# esphome-espnow-tree-ha/rootfs/etc/services.d/esphome-espnow-tree-ha/run
#!/usr/bin/with-contenv sh
set -eu

cd /opt/espnow-tree
exec uvicorn app.main:app --host 0.0.0.0 --port 8099 --timeout-graceful-shutdown 3
